# server/ansible/roles/app/templates/supervisor.conf.j2
[program:breta]
command = /www/backend/breta/gunicorn_start.sh
user = root
redirect_stderr = true
